# cloudgaming-platform/Cloudgaming-platform-main/cloud_gaming_system/vm-host-linux/vm-host.service
[Unit]
Description=Cloud Gaming VM Host Capture Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cloud-gaming-vm-host
EnvironmentFile=-/opt/cloud-gaming-vm-host/.env
Environment=DISPLAY=:0
ExecStart=/opt/cloud-gaming-vm-host/venv/bin/python3 -m src.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=false
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
